# docker/crontab/tasks
# Sincronizacao principal
17 1       * * *   www-data python /opt/suap/manage.py sync_suap >> /opt/suap/deploy/logs/cron/sync_suap.log 2>&1

# ldap
40 8-22    * * *   www-data python /opt/suap/manage.py sync_ldap_parcial >>  /opt/suap/deploy/logs/cron/sync_ldap_parcial.log
0 23       * * *   www-data python /opt/suap/manage.py sync_ldap_full >>  /opt/suap/deploy/logs/cron/sync_ldap_full.log

# acompanhamentos catalogo provedor de servico
22 *       * * *   www-data timeout 1200 python /opt/suap/manage.py enviar_email_avaliacao_de_servicos >>  /opt/suap/deploy/logs/cron/acompanhamentos.log 2>&1

# processando pendencias do barramento
12 *       * * *   www-data python /opt/suap/manage.py processar_pendencias_barramento >> /opt/suap/deploy/logs/cron/barramento.log 2>&1

# importar siape via ws
7 21       * * *   www-data python /opt/suap/manage.py atualizar_servidor_completo_ws --completo >>  /opt/suap/deploy/logs/cron/servidores_ws.log 2>&1

# gerar notificações para usuários
27 11,16   * * *   www-data python /opt/suap/manage.py enviar_notificacoes_notifica_govbr >>  /opt/suap/deploy/logs/cron/notificacoes.txt

# cnpq
1 19       * * *   www-data python /opt/suap/manage.py cnpq_importar --xml >> /opt/suap/deploy/logs/cron/cnpq_importar.log 2>&1
10 *       * * *   www-data python /opt/suap/manage.py atualizar_grupos_pesquisa_dos_servidores  >> /opt/suap/deploy/logs/cron/atualizar_grupos_pesquisa_dos_servidores.log 2>&1

# ultima execução cron
* *        * * *   www-data echo $(date) > /opt/suap/deploy/logs/cron/ultimo.log

# tmpreaper
0  *       * * *   www-data tmpreaper 4h /tmp --showdeleted >>  /opt/suap/deploy/logs/cron/tmpreaper.log

# correcao automatica de data/hora
59 *       * * *   www-data /usr/sbin/ntpdate a.ntp.br >> /opt/suap/deploy/logs/cron/ntp.log

# validando passaportes vacinais de COVID
10 *       * * *   www-data python /opt/suap/manage.py atualiza_validade_passaportes_vacinais >> /opt/suap/deploy/logs/cron/validade_passaportes_vacinais.log 2>&1
10 22      * * *   www-data python /opt/suap/manage.py atualiza_passaportes_vacinais_pendentes >> /opt/suap/deploy/logs/cron/atualiza_passaportes_vacinais_pendentes.log 2>&1

# sincroniza as assinaturas digitais de diplomas com a ferramenta RAP da RNP
0,15,30,45 * * * * www-data python /opt/suap/manage.py edu_sincronizar_assinatura_digital >> /opt/suap/deploy/logs/cron/sync_diploma_digital.log 2>&1
